# usb_modeswitch/92-ifupdown.rules
SUBSYSTEM=="net", ACTION=="add", RUN+="ifup %E{INTERFACE}" 
SUBSYSTEM=="net", ACTION=="remove", RUN+="ifdown %E{INTERFACE}"

