# ltib/skeleton/rootfs/etc/inittab
#
# inittab       This file describes how the INIT process should set up
#               the system in a certain run-level.
#		Please read a standard system's man pages for inittab
#

# Runlevels. Only these runlevels used by MBARI embedded systems:
#   3 - (default) multiuser mode
#   6 - reboot

# Default Runlevel to Go into
dflt:3:initdefault:

# System Startup. Do these process before anything else
strt::sysinit:/etc/rc.d/sysinit

# Process to do for each specific Runlevel
r3:3:wait:/etc/rc.d/rc 3
r6:6:wait:/etc/rc.d/rc 6

# Start console logins
co:12345:respawn:/sbin/getty -L 115200 ttyS0

# Keep basic network services running
net:12345:respawn:/usr/sbin/inetd -f

# Keep the freewave ppp interface up
#ppp0:12345:respawn:/usr/sbin/pppd nodetach call shore unit 0 logfile /dev/null

# Keep the freewave slip interface up
#sl0:12345:respawn:/sbin/slattach -P /var/run/sl0.pid -i "/usr/sbin/slinit sl0" -lp slip -s 115200 /dev/freewave

# Keep up VPN link to MBARI (after allowing time for network to start)
#vpn:12345:respawn:/bin/sh -c "sleep 300; exec /usr/sbin/inadyn"

# Keep ssh tunnel to MBARI up
#tun:12345:respawn:/usr/sbin/tunnel2shore

# keep an espclient up on ttyS3
#esp:3:respawn:/sbin/getty -L 9600 -nl /etc/esp/client ttyS3
