# ltib/skeleton/rootfs/etc/inittab
#
# inittab       This file describes how the INIT process should set up
#               the system in a certain run-level.
#		Please read a standard system's man pages for inittab
#

# Runlevels. Only these runlevels used by MBARI embedded systems:
#   3 - (default) multiuser mode
#   6 - reboot

# System Startup. Do these process before anything else
strt::sysinit:/etc/rc.d/sysinit

# Default Runlevel to Go into
dflt:3:initdefault:

# Process to do for each specific Runlevel
r3:3:wait:/etc/rc.d/rc 3
r6:6:wait:/etc/rc.d/rc 6

# Start console logins
co:12345:respawn:/sbin/getty -L 115200 ttyS0

#buoy:12345:respawn:/sbin/getty -L 115200 buoyConsole

#bring up deferred net interfaces
ifup:3:once:/lib/udev/catchup

# Keep basic network services running
net:12345:respawn:/usr/sbin/inetd -f

# start handling device events (optionally after specified delay in seconds)
#DEVD:3:once:/sbin/startudevd 90

# Keep the freewave ppp interface up
#ppp0:12345:respawn:/usr/sbin/pppd nodetach call shore unit 0 logfile /dev/null

# Keep the freewave slip interface up
#sl0:12345:respawn:/sbin/slattach -P /var/run/sl0.pid -i "/usr/sbin/slinit sl0" -lp slip -s 115200 /dev/freewave

# Keep up VPN link to MBARI (after allowing time for network to start)
#vpn:12345:respawn:/bin/sh -c "sleep 300; exec /usr/sbin/inadyn"

# reset USB modem when Wireguard session fails (at most 1 of the examples below)
#wg:12345:respawn:/usr/sbin/wgwatch wg2shore /etc/ppp/reset-modem
# reset serial modem [vzlte] when Wireguard session fails
#wg:12345:respawn:/usr/sbin/wgwatch wg2shore "ifdown vzlte; noUSBreset=1 /etc/ppp/reset-modem; sleep 30; ifup vzlte"
# special case for dealing with AT&T's OneRate LTM1 service on wwan0
#wg:12345:respawn:/usr/sbin/wgwatch wg2shore "ifdown wwan0; sleep 95; resetUSB"

# Keep ssh tunnel to MBARI up
#tun:12345:respawn:/usr/sbin/tunnel2shore

# keep an espclient up on ttyS3
#esp:3:respawn:/sbin/getty -L 9600 -nl /etc/esp/client ttyS3

# call back if modem RINGs
#ring:3:respawn:/usr/sbin/ringBack

#every hour, cycle USB and modem power if modem unresponsive
#mind:3:respawn:/bin/sh -c "sleep 1h && detachModem iridium 2>/dev/null"

# start ESP server running yourMission script
#ESP:3:once:/bin/su - esp -c "sleep 5 && ESPcmdPort=9999 start esp yourMission"
