# ltib/skeleton/rootfs/etc/rc.local
#!/bin/sh
#enable tracing specified scripts
#mkdir /run/trace && cd /run/trace && touch rerouteESPshore

#bring up satellite link if it is already powered
#power sat | grep "ON$" && ifup certus
