# ltib/skeleton/rootfs/etc/rc.local
#!/bin/sh

#enable tracing the specified scripts
#mkdir /run/trace && cd /run/trace && touch rerouteESPshore wgwatch
