# usbmode/40-usbmode.rules
# Derived from usb-modeswitch-data, version 20170805
#
SUBSYSTEMS=="usb", ACTION=="add", RUN+="usbmode -s"
