# ltib/skeleton/rootfs/etc/udev/rules.d/92-ifupdown.rules
#ignore virtual network interfaces
SUBSYSTEM!="net", GOTO="notNetIf"
ENV{DEVPATH}=="/devices/virtual/*", GOTO="notTTYIf"

ACTION=="add", RUN+="ifup %E{INTERFACE}"
ACTION=="remove", RUN+="ifdown %E{INTERFACE}"

LABEL="notNetIf"
SUBSYSTEM!="tty", GOTO="notTTYif"

### handle serial USB modems
ACTION=="remove", RUN+="ifTTYdown"
ACTION!="add", GOTO="updownEND"

#Quectel BG96
ATTRS{idVendor}=="2c7c", ATTRS{idProduct}=="0296", ENV{ID_IFACE}=="03",\
  SYMLINK+="AT/bg96", GROUP="root", RUN+="ifup bg96"

#Iridium Satphone Handset
ENV{ID_MODEL}=="Iridium_*_Handset", ENV{ID_IFACE}=="00",\
  SYMLINK+="AT/iridium", GROUP="root", RUN+="ifup iridium"

#delay SARA-R4 ifup by after 4 seconds to try to avoid crashing its firmware
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="90b2", ENV{ID_IFACE}=="02",\
  SYMLINK+="AT/saraR4", GROUP="root", RUN+="ifup saraR4 4"

GOTO="updownEND"
LABEL="notTTYif"

#cycle USB power whenever the SARA-R4 modem fireware crashes
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="90b2",\
  ATTRS{product}=="*USB__BULK", RUN+="/usr/sbin/resetUSB 10", GOTO="notTTYIf"

LABEL="updownEND"
