# ltib/skeleton/rootfs/etc/udev/rules.d/92-ifupdown.rules
#12/2/25 brent@mbari.org
#ignore virtual network interfaces
SUBSYSTEM!="net", GOTO="notNetIf"
ENV{DEVPATH}=="/devices/virtual/*", GOTO="updownEND"
KERNEL=="mon.*", GOTO="updownEND"

ACTION=="remove", RUN+="ifdown %k", GOTO="updownEND"
ACTION!="add", GOTO="updownEND"

ATTRS{busnum}=="", GOTO="UPnow"

#allow extra time for Telit to boot
ATTRS{idVendor}=="1bc7", ENV{delay}="30"

#ensure stickup runs before attempting to bring up SIMCom A76xx
ATTRS{idVendor}=="1e0e", ENV{delay}="8"

#Pantech does not require much delay
ATTRS{idVendor}=="10a9", ENV{delay}="1"

RUN+="ifLink USB%s{busnum}:%s{devnum} %k"

LABEL="UPnow"
RUN+="ifup %k %E{delay}"
GOTO="updownEND"

LABEL="notNetIf"
SUBSYSTEM!="tty", GOTO="updownEND"

### symlink AT ports on recognized USB modems to /dev/AT/*modem net iface*
ACTION!="add", GOTO="updownEND"

ATTRS{idVendor}=="1e0e", GOTO="Arcom"
ATTRS{idVendor}=="05c6", GOTO="SIMCom"
ATTRS{idVendor}=="1bc7", GOTO="Telit"
ATTRS{idVendor}=="10a9", GOTO="Pantech"
ATTRS{idVendor}=="2c7c", GOTO="Quectel"
ATTRS{idVendor}=="12d1", GOTO="Huawei"

#Iridium Satphone Handset
ATTRS{idVendor}=="1edd",ENV{ID_MODEL}=="Iridium_*_Handset",ENV{ID_IFACE}=="00",\
  GROUP:="root", SYMLINK+="AT/iridium", RUN+="ifup iridium"
GOTO="updownEND"


LABEL="Arcom"
ATTRS{idProduct}=="9011", GOTO="A76XX"
ATTRS{idProduct}=="9205", GOTO="A7080"
GOTO="updownEND"

LABEL="A7080"
ENV{ID_IFACE}=="00", GOTO="diagPort"
ENV{ID_IFACE}=="02", GOTO="ATport"
ENV{ID_IFACE}=="03", GOTO="PPPport"
ENV{ID_IFACE}=="01", GOTO="GPSport"
GOTO="updownEND"

LABEL="A76XX"
ENV{ID_IFACE}=="02", GOTO="diagPort"
ENV{ID_IFACE}=="04",\
  RUN+="stickup $env{DEVNAME} 'AT+USBNETIP=1'", GOTO="ATport"
ENV{ID_IFACE}=="05", GOTO="PPPport"
ENV{ID_IFACE}=="06", GOTO="GPSport"
GOTO="updownEND"


LABEL="SIMCom"
ATTRS{idProduct}=="9330", GOTO="SIM76XX"
#delay SARA-R4 ifup by a few seconds
ATTRS{idProduct}=="90b2", ENV{ID_IFACE}=="02",\
  GROUP:="root", SYMLINK+="AT/saraR4", RUN+="ifup saraR4 14"
GOTO="updownEND"

LABEL="SIM76XX"
ENV{ID_IFACE}=="02", GOTO="ATport"
ENV{ID_IFACE}=="04", GOTO="diagPort"
ENV{ID_IFACE}=="06", GOTO="PPPport"
ENV{ID_IFACE}=="08", GOTO="GPSport"
GOTO="updownEND"


#Telit ME910 (prod 1102 is a USB CDCether device)
LABEL="Telit"
ATTRS{idProduct}=="1102", ENV{ID_IFACE}=="02", GOTO="ATport"
ATTRS{idProduct}=="1101",\
  GROUP:="root", SYMLINK+="AT/me910", RUN+="ifup me910 15"
GOTO="updownEND"


#Pantech UML295
LABEL="Pantech"
ATTRS{idProduct}=="6064", ENV{ID_IFACE}=="06", GOTO="ATport"
ATTRS{idProduct}=="6064", ENV{ID_IFACE}=="04", GOTO="GPSport"
GOTO="updownEND"


LABEL="Quectel"
ATTRS{idProduct}=="0901", GOTO="EG912"
ATTRS{idProduct}=="0296",ENV{ID_IFACE}=="03",\
  GROUP:="root",SYMLINK+="AT/bg96",RUN+="ifup bg96 10"
GOTO="updownEND"

LABEL="EG912"
ENV{ID_IFACE}=="02", GOTO="ATport"
ENV{ID_IFACE}=="07", GOTO="PPPport"
ENV{ID_IFACE}=="08", GOTO="GPSport"
GOTO="updownEND"


#Huawei modems in "stick" mode need to receive AT cmd to attach to the network
LABEL="Huawei"
ATTRS{idProduct}=="1506", ENV{ID_IFACE}=="01", GOTO="ATport"
GOTO="updownEND"


LABEL="GPSport"
MODE="0640", SYMLINK+="AT/gps", ENV{newPort}="AT/#-gps", GOTO="addPort"

LABEL="diagPort"
ENV{newPort}="AT/#-diag", GOTO="rootPort"

LABEL="PPPport"
ENV{newPort}="AT/#-ppp", GOTO="rootPort"

LABEL="ATport"
ENV{newPort}="AT/#"
LABEL="rootPort"
GROUP:="root"
LABEL="addPort"
RUN+="ifTTY USB%s{busnum}:%s{devnum} %E{newPort}"

LABEL="updownEND"
