# ltib/skeleton/rootfs/etc/udev/rules.d/93-ttylogin.rules
#allow one login from newly inserted USB serial port
ACTION=="add", SUBSYSTEM=="tty", DRIVERS=="pl2303|ftdi_sio|mct_u232", \
  RUN+="ttyLogin"
